# fabric/systemd/forenshield-fabric-gateway.service
[Unit]
Description=ForenShield Fabric Anchor Gateway
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/forenshield-infra/fabric/gateway
EnvironmentFile=/home/ubuntu/forenshield-infra/fabric/gateway/.env
ExecStart=/usr/bin/npm start
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
