# fabric/systemd/forenshield-fabric-gateway.service
[Unit]
Description=ForenShield Fabric Anchor Gateway
After=forenshield-fabric-network.service docker.service network-online.target
Wants=forenshield-fabric-network.service docker.service
Requires=docker.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/forenshield-infra/fabric/gateway
EnvironmentFile=/home/ubuntu/forenshield-infra/fabric/gateway/.env
ExecStart=/usr/bin/npm start
Restart=on-failure
RestartSec=15
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
